# excludetiktok
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org

START=90
INTERFACE=br-lan
PORT=1337

# check if configuration exists
[ -e "/etc/redsocks.conf" ] || exit 0

# List of Bangladeshi IP address ranges to exclude from redsocks
BD_IPS="\
101.33.20.236/32 \
101.33.21.149/32 \
101.33.26.117/32 \
101.33.26.118/32 \
101.33.26.135/32 \
101.33.26.137/32 \
101.33.26.17/32 \
101.33.26.196/32 \
101.33.26.34/32 \
101.33.26.41/32 \
101.33.27.105/32 \
101.33.27.110/32 \
101.33.27.120/32 \
101.33.27.58/32 \
101.33.27.61/32 \
103.12.74.137/32 \
103.12.74.138/32 \
103.197.254.228/32 \
103.197.254.238/32 \
104.251.234.44/32 \
109.61.80.118/32 \
109.61.80.123/32 \
109.61.80.76/32 \
115.108.61.19/32 \
115.108.61.20/32 \
115.108.61.21/32 \
115.108.61.23/32 \
115.108.61.44/32 \
129.227.80.18/32 \
143.244.51.17/32 \
143.244.51.18/32 \
143.42.210.140/32 \
143.42.210.83/32 \
143.42.210.87/32 \
143.42.210.91/32 \
143.42.210.92/32 \
146.75.118.73/32 \
147.160.178.35/32 \
147.160.179.32/32 \
147.160.180.42/32 \
148.153.241.20/32 \
148.153.241.23/32 \
148.153.241.40/32 \
149.102.250.150/32 \
149.102.250.157/32 \
149.102.250.172/32 \
149.102.250.173/32 \
149.102.250.174/32 \
150.109.90.4/32 \
151.101.38.73/32 \
154.85.73.249/32 \
154.85.73.250/32 \
154.85.73.251/32 \
154.85.73.252/32 \
154.85.73.253/32 \
156.225.96.118/32 \
156.225.96.120/32 \
156.225.96.122/32 \
156.225.96.123/32 \
156.225.96.200/32 \
157.185.143.23/32 \
199.232.170.73/32 \
207.211.211.105/32 \
207.211.211.113/32 \
207.211.211.114/32 \
207.211.211.123/32 \
207.211.211.226/32 \
212.102.56.122/32 \
212.102.56.73/32 \
212.102.56.82/32 \
212.102.56.87/32 \
212.102.56.89/32 \
23.210.93.145/32 \
23.210.93.146/32 \
23.210.93.147/32 \
23.210.93.187/32 \
23.210.93.193/32 \
23.210.93.210/32 \
23.210.93.217/32 \
23.210.93.225/32 \
23.210.93.226/32 \
23.212.164.10/32 \
23.212.164.11/32 \
23.212.164.113/32 \
23.212.164.120/32 \
23.212.164.121/32 \
23.212.164.122/32 \
23.212.164.123/32 \
23.212.164.129/32 \
23.212.164.131/32 \
23.212.164.139/32 \
23.212.164.144/32 \
23.212.164.146/32 \
23.212.164.16/32 \
23.212.164.17/32 \
23.212.164.18/32 \
23.212.164.19/32 \
23.212.164.24/32 \
23.212.164.25/32 \
23.212.164.26/32 \
23.212.164.27/32 \
23.212.164.32/32 \
23.212.164.33/32 \
23.212.164.34/32 \
23.212.164.35/32 \
23.212.164.40/32 \
23.212.164.41/32 \
23.212.164.42/32 \
23.212.164.43/32 \
23.212.164.48/32 \
23.212.164.49/32 \
23.212.164.50/32 \
23.212.164.51/32 \
23.212.164.56/32 \
23.212.164.57/32 \
23.212.164.58/32 \
23.212.164.59/32 \
23.212.164.64/32 \
23.212.164.65/32 \
23.212.164.8/32 \
23.212.164.9/32 \
23.217.118.166/32 \
23.217.118.167/32 \
23.57.76.22/32 \
23.57.76.23/32 \
23.57.76.24/32 \
23.57.76.37/32 \
23.57.76.44/32 \
23.57.76.45/32 \
23.57.76.46/32 \
23.57.76.47/32 \
23.57.76.52/32 \
23.57.76.53/32 \
23.57.76.54/32 \
23.57.76.55/32 \
23.57.76.60/32 \
23.57.76.61/32 \
23.57.76.62/32 \
23.57.76.63/32 \
23.57.76.69/32 \
23.58.120.105/32 \
23.58.120.106/32 \
23.58.120.112/32 \
23.58.120.113/32 \
23.58.120.129/32 \
23.58.120.136/32 \
23.58.120.138/32 \
23.58.120.155/32 \
23.58.120.160/32 \
23.58.120.162/32 \
23.58.120.89/32 \
23.63.109.10/32 \
23.63.109.104/32 \
23.63.109.106/32 \
23.63.109.11/32 \
23.63.109.114/32 \
23.63.109.129/32 \
23.63.109.136/32 \
23.63.109.145/32 \
23.63.109.147/32 \
23.63.109.16/32 \
23.63.109.160/32 \
23.63.109.168/32 \
23.63.109.170/32 \
23.63.109.176/32 \
23.63.109.18/32 \
23.63.109.19/32 \
23.63.109.219/32 \
23.63.109.25/32 \
23.63.109.27/32 \
23.63.109.32/32 \
23.63.109.33/32 \
23.63.109.35/32 \
23.63.109.42/32 \
23.63.109.48/32 \
23.63.109.51/32 \
23.63.109.56/32 \
23.63.109.57/32 \
23.63.109.59/32 \
23.63.109.73/32 \
23.63.109.8/32 \
23.63.109.82/32 \
23.63.109.83/32 \
23.63.109.90/32 \
23.63.109.97/32 \
23.63.109.98/32 \
23.63.111.146/32 \
23.63.111.153/32 \
23.63.111.162/32 \
23.63.111.170/32 \
23.63.111.176/32 \
23.63.111.177/32 \
23.63.111.192/32 \
23.63.111.194/32 \
23.63.111.203/32 \
23.72.88.16/32 \
23.72.88.35/32 \
23.72.88.50/32 \
23.72.88.57/32 \
23.72.88.73/32 \
23.72.88.88/32 \
23.72.88.96/32 \
23.72.88.97/32 \
31.3.2.61/32 \
31.3.2.64/32 \
31.3.2.67/32 \
34.102.147.19/32 \
34.102.164.249/32 \
34.107.238.235/32 \
34.110.200.237/32 \
34.36.65.236/32 \
34.36.71.42/32 \
34.36.75.57/32 \
34.36.80.120/32 \
34.96.106.127/32 \
38.60.148.168/32 \
38.60.148.170/32 \
38.60.148.187/32 \
38.60.148.188/32 \
43.132.72.158/32 \
43.132.73.69/32 \
43.132.73.71/32 \
43.132.73.73/32 \
43.132.73.82/32 \
43.132.80.102/32 \
43.132.80.128/32 \
43.132.80.131/32 \
43.132.80.182/32 \
43.132.80.207/32 \
43.132.80.228/32 \
43.132.80.241/32 \
43.132.80.66/32 \
43.132.81.144/32 \
43.132.81.162/32 \
43.132.81.172/32 \
43.132.81.182/32 \
43.132.81.91/32 \
43.132.84.154/32 \
43.152.142.241/32 \
43.152.23.169/32 \
43.152.24.221/32 \
43.152.24.61/32 \
43.152.25.73/32 \
43.159.70.154/32 \
43.175.34.28/32 \
43.175.34.30/32 \
43.175.44.149/32 \
43.175.45.213/32 \
43.175.45.242/32 \
43.175.46.122/32 \
43.175.46.123/32 \
43.175.46.241/32 \
43.175.46.242/32 \
43.175.46.36/32 \
43.175.46.74/32 \
43.175.47.145/32 \
43.175.47.35/32 \
43.175.47.46/32 \
43.175.47.51/32 \
43.175.48.11/32 \
43.175.48.206/32 \
43.175.48.237/32 \
43.175.48.254/32 \
43.175.48.68/32 \
43.175.48.83/32 \
71.18.1.248/32 \
71.18.255.144/32 \
71.18.255.145/32 \
71.18.255.146/32 \
71.18.255.147/32 \
71.18.39.231/32 \
71.18.39.232/32 \
71.18.71.119/32 \
71.18.71.129/32 \
79.127.213.226/32 \
79.127.213.229/32 \
79.127.213.230/32 \
79.127.235.19/32 \
79.127.235.23/32 \
84.17.61.17/32 \
84.17.61.18/32 \
89.187.163.236/32 \
98.98.121.32/32 \
98.98.121.43/32 \
98.98.227.10/32 \
98.98.227.30/32 \
98.98.227.46/32 \
98.98.227.47/32 "

iptable_start() {
    /bin/echo -n "running redsocks iptables ..."

    # Run iptable commands
    iptables -t nat -N REDSOCKS

    # Exclude private networks and Bangladeshi IP addresses
    iptables -t nat -A REDSOCKS -d 0.0.0.0/8 -j RETURN
    iptables -t nat -A REDSOCKS -d 10.0.0.0/8 -j RETURN
    iptables -t nat -A REDSOCKS -d 127.0.0.0/8 -j RETURN
    iptables -t nat -A REDSOCKS -d 169.254.0.0/16 -j RETURN
    iptables -t nat -A REDSOCKS -d 172.16.0.0/12 -j RETURN
    iptables -t nat -A REDSOCKS -d 192.168.0.0/16 -j RETURN
    iptables -t nat -A REDSOCKS -d 224.0.0.0/4 -j RETURN
    iptables -t nat -A REDSOCKS -d 240.0.0.0/4 -j RETURN

    for IP in $BD_IPS; do
        iptables -t nat -A REDSOCKS -d $IP -j RETURN
    done

    iptables -t nat -A REDSOCKS -p tcp -j REDIRECT --to-ports ${PORT}

    iptables -t nat -A PREROUTING -i ${INTERFACE} -p tcp -j REDSOCKS

    iptables -A INPUT -i br-lan -p tcp --dport ${PORT} -j ACCEPT

    /bin/echo " done"
}

iptable_stop() {
    /bin/echo -n "cleaning redsocks iptables ..."

    # Run iptable commands
    iptables -t nat -F REDSOCKS
    iptables -t nat -F PREROUTING
    iptables -t nat -F POSTROUTING
    iptables -F INPUT
    iptables -F FORWARD
    iptables -t nat -X REDSOCKS

    /bin/echo " done"
}

start() {
    if [ -e "/var/run/redsocks.pid" ]; then
        echo "redsocks is already running"
        exit 0
    fi

    /bin/echo -n "running redsocks ..."

    # startup the safety-wrapper for the daemon
    /usr/sbin/redsocks -p /var/run/redsocks.pid

    /bin/echo " done"
    iptable_start
}

stop() {
    if [ ! -e "/var/run/redsocks.pid" ]; then
        echo "redsocks is not running"
        exit 0
    fi

    /bin/echo -n "stopping redsocks ..."

    # kill the process
    /bin/kill $(cat /var/run/redsocks.pid)
    rm /var/run/redsocks.pid

    echo " done"
    iptable_stop

    /bin/echo -n "restarting firewall ..."
    /etc/init.d/firewall restart &> /dev/null
    /bin/echo " done"
}
